# dockerfiles/supervisord.conf
[supervisord]
nodaemon=true

[program:waggle-dance]
command=/opt/waggle-dance/bin/waggle-dance.sh start
startsecs=0
autorestart=false
startretries=1
autostart=true